# config/hosthive-api.service
[Unit]
Description=HostHive API Server (uvicorn)
Documentation=https://hosthive.dev
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=hosthive
Group=hosthive
WorkingDirectory=/opt/hosthive
EnvironmentFile=/opt/hosthive/config/secrets.env

ExecStart=/opt/hosthive/venv/bin/uvicorn \
    app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info \
    --access-log

Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hosthive /var/log/hosthive
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hosthive-api

[Install]
WantedBy=multi-user.target
