# config/hosthive-worker.service
[Unit]
Description=HostHive Celery Worker (background tasks)
Documentation=https://hosthive.dev
After=network.target redis-server.service postgresql.service
Wants=redis-server.service postgresql.service

[Service]
Type=simple
User=hosthive
Group=hosthive
WorkingDirectory=/opt/hosthive
EnvironmentFile=/opt/hosthive/config/secrets.env

ExecStart=/opt/hosthive/venv/bin/celery \
    -A app.worker \
    worker \
    --loglevel=info \
    --concurrency=2 \
    --max-tasks-per-child=100 \
    -Q default,emails,backups,dns

Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hosthive /var/log/hosthive /tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hosthive-worker

[Install]
WantedBy=multi-user.target
