# config/novapanel-worker.service
[Unit]
Description=HostHive Celery Worker (background tasks)
Documentation=https://hosthive.dev
After=network.target redis-server.service postgresql.service
Wants=redis-server.service postgresql.service

[Service]
Type=simple
User=novapanel
Group=novapanel
WorkingDirectory=/opt/novapanel
EnvironmentFile=/opt/novapanel/config/secrets.env

ExecStart=/opt/novapanel/venv/bin/celery \
    -A app.worker \
    worker \
    --loglevel=info \
    --concurrency=2 \
    --max-tasks-per-child=100 \
    -Q default,emails,backups,dns

Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/novapanel /var/log/novapanel /tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=novapanel-worker

[Install]
WantedBy=multi-user.target
